# monitoring.service
[Unit]
Description=Laravel Monitoring Stack (Loki, Prometheus, Grafana)
Requires=docker.service
After=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/var/www/prohelper
ExecStart=/bin/bash -c 'cd /var/www/prohelper && (command -v docker-compose >/dev/null && docker-compose up -d || docker compose up -d)'
ExecStop=/bin/bash -c 'cd /var/www/prohelper && (command -v docker-compose >/dev/null && docker-compose down || docker compose down)'
ExecReload=/bin/bash -c 'cd /var/www/prohelper && (command -v docker-compose >/dev/null && docker-compose restart || docker compose restart)'
TimeoutStartSec=0
TimeoutStopSec=30
Restart=on-failure
RestartSec=5

# Безопасность
User=www-data
Group=www-data

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitoring-stack

[Install]
WantedBy=multi-user.target